# .config/systemd/user/bugwarrior-pull.service
[Unit]
Description=Run bugwarrior-pull to fetch tasks

[Service]
ExecStart=/usr/bin/bugwarrior pull
Type=oneshot

[Install]
WantedBy=default.target
