# .config/systemd/user/newsboatsync.service
[Unit]
Description=Newsboat syncronization service

[Service]
Type=oneshot
ExecStart=%h/.config/scripts/newsboatsync
